# policy-updater/policies-setup.service
[Unit]
Description=Setup policies
After=init-done.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/policy-updater
